# deklan-fusion/etc/systemd/fusion-monitor.service
[Unit]
Description=Deklan Fusion Monitor Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/deklan-fusion
Environment="PATH=/opt/deklan-fusion/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/opt/deklan-fusion/venv/bin/python3 -m monitor.monitor
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target


