# etc/systemd/system/deklan-monitor.service
[Unit]
Description=Deklan Fusion Node Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/deklan-fusion
EnvironmentFile=/opt/deklan-fusion/.env
ExecStart=/opt/deklan-fusion/.venv/bin/python /opt/deklan-fusion/monitor/monitor.py

# Logging
StandardOutput=journal
StandardError=journal
